# .config/autostart/dock1.desktop
[Desktop Entry]
Type=Application
Exec=plank -n dock1
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Dock 1