# .config/autostart/setup-monitors.desktop
[Desktop Entry]
Type=Application
Exec=./setup-monitors.sh
X-GNOME-Autostart-enabled=true
Name=Setup Monitors
